# deployment/chainmaker-dev.service
[Unit]
Description=ChainMaker Dev Environment
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/chainmaker-dev
EnvironmentFile=/var/www/chainmaker-dev/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chainmaker-dev

# Simplified security settings to avoid NAMESPACE errors
PrivateTmp=true
ReadWritePaths=/var/www/chainmaker-dev/public/user_images

[Install]
WantedBy=multi-user.target
