# deployment/chainmaker.service
[Unit]
Description=ChainMaker Web Application
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/chainmaker
EnvironmentFile=/var/www/chainmaker/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chainmaker

# Simplified security settings to avoid NAMESPACE errors
PrivateTmp=true
ReadWritePaths=/var/www/chainmaker/public/user_images

[Install]
WantedBy=multi-user.target